# scripts/companion/systemd/q6a-announce.service
[Unit]
Description=Companion object announcer: /vision/detections -> spoken /robot/speak ("I see a chair")
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=radxa
EnvironmentFile=-/etc/default/ippolit-robot
ExecStart=/bin/bash -c "source /opt/ros/jazzy/setup.bash && exec python3 /home/radxa/ros/q6a_announce.py"
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
