# scripts/companion/systemd/q6a-slam-toolbox.service
[Unit]
Description=Companion slam_toolbox (async): /scan + odom->base_link -> map->odom TF + /map occupancy grid
# q6a_laser_odom now lives in ippolit-core.service (ippolit_localization pkg) since 2026-07-12's
# A1 migration -- After= updated accordingly (was q6a-laser-odom.service, since retired).
After=network-online.target ippolit-core.service
Wants=network-online.target

[Service]
KillSignal=SIGINT
TimeoutStopSec=8
Type=simple
User=radxa
EnvironmentFile=-/etc/default/ippolit-robot
ExecStart=/bin/bash -c "source /opt/ros/jazzy/setup.bash && exec ros2 run slam_toolbox async_slam_toolbox_node --ros-args --params-file /home/radxa/ros/slam_toolbox.yaml"
# Jazzy slam_toolbox is a LIFECYCLE node: without configure+activate it idles in 'unconfigured'
# (no solver, no /scan subscription, no errors). ExecStartPost walks it to 'active'.
ExecStartPost=/home/radxa/ros/slam_lifecycle_up.sh
# Map save/resume: see q6a-map-persist.service (a dedicated node, not a hook on this unit).
TimeoutStartSec=180
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
