# scripts/companion/systemd/q6a-llmd.service
[Unit]
Description=Q6A resident Genie LLM daemon (Llama 3.2 1B on Hexagon NPU)
After=llama-prewarm.service
Wants=llama-prewarm.service
# Anti-restart-storm: the daemon self-exits on a dead cDSP context (SSR) so systemd reloads it, but a
# genuinely-wedged cDSP must NOT be restart-stormed (that needs a reboot to clear). Cap to 4 starts / 5 min;
# past that systemd stops trying and leaves it failed (a human/monitor intervenes) rather than storming.
StartLimitIntervalSec=300
StartLimitBurst=4

[Service]
KillSignal=SIGINT
TimeoutStopSec=8
Type=simple
User=radxa
SupplementaryGroups=render
WorkingDirectory=/home/radxa/llama-1b
Environment=LD_LIBRARY_PATH=/home/radxa/llama-1b
Environment=ADSP_LIBRARY_PATH=/home/radxa/llama-1b
Environment=Q6A_LLM_SOCK=/tmp/q6a-llm.sock
ExecStart=/usr/bin/python3 /home/radxa/llama-1b/q6a_llmd.py
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
