# ros2_ws/src/ippolit_bringup/systemd/ippolit-viz.service
[Unit]
Description=IPPOLIT visualization: foxglove_bridge websocket (restart blast-radius group 4 of 4)
After=network-online.target ippolit-core.service
Wants=network-online.target

[Service]
KillSignal=SIGINT
TimeoutStopSec=8
Type=simple
User=radxa
EnvironmentFile=-/etc/default/ippolit-robot
ExecStart=/bin/bash -c "source /opt/ros/jazzy/setup.bash && source ~/ippolit/ros2_ws/install/setup.bash && exec ros2 launch ippolit_bringup viz.launch.xml"
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
